# common/etc/update-motd.d/00-header
#!/bin/sh
printf "Welcome to %s (%s %s %s)\n" "Cubian" "$(uname -o)" "$(uname -r)" "$(uname -m)" > /var/run/motd.dynamic
