# common/etc/update-motd.d/10-help-text
#!/bin/sh
cat >>/var/run/motd.dynamic<<EOF

 * Documentation:  http://cubian.org
EOF
