# common/etc/update-motd.d/20-package-info
#!/bin/sh

N=/var/run/motd.package

echo "/usr/lib/cubian/cubian-motd-package" | at now

if [ -f $N ];then
echo $N >> /var/run/motd.dynamic
fi
