# common/etc/network/if-up.d/show-morsecode-ip
#!/bin/sh

set -e

PATH=/sbin:/bin:/usr/sbin:/usr/bin

if [ -x /usr/lib/cubian/cubian-networklight ];then
	echo "/usr/lib/cubian/cubian-networklight 10" | at now
fi
exit 0
